# files/etc/uci-defaults/disable-dnsmasq.sh
#!/bin/sh
echo "Disabling dnsmasq..."
/etc/init.d/dnsmasq disable
